# supervisored_services.ini
[program:mosquitto]
command=/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
autorestart=true
numprocs=1
[program:shelly-handler]
command=/opt/shelly/src/entrypoint.sh -c /opt/configs/tdm.conf
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
autorestart=true
numprocs=1
